# roles/podman/files/ente/ente-museum.container
[Unit]
Description=Ente Museum
After=ente-postgres.service
After=ente-garage.service
Wants=ente-postgres.service
Wants=ente-garage.service

[Container]
ContainerName=ente-museum
Image=ghcr.io/ente-io/server:6f260f1f7de6f0134b107677a7a40c4720ed05b3@sha256:e1b615b261391b79a9d27f00783f70b909a18c67960cc038846a9e928fe79903
AutoUpdate=registry

Network=nginx-ente.network
Network=ente-backend.network
NetworkAlias=ente-museum

Volume=ente-config:/config:ro,Z
ReadOnly=true
Mount=type=tmpfs,tmpfs-size=256M,tmpfs-mode=1777,U,noexec,nosuid,nodev,destination=/tmp

Environment=ENTE_CREDENTIALS_FILE=/config/credentials.yaml

User=300015
Group=300015

GlobalArgs=--runtime /usr/local/bin/runsc
SecurityLabelDisable=true
NoNewPrivileges=true
DropCapability=ALL

[Service]
LockPersonality=true
ProtectClock=true
RestrictRealtime=true
RuntimeDirectoryMode=700
SystemCallArchitectures=native
SystemCallFilter=~@obsolete
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
