# roles/podman/templates/haproxy/haproxy-rotate-tls-ticket-keys.service.j2
[Unit]
Description=Rotate HAProxy TLS ticket keys
After=haproxy.service
Wants=haproxy.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/haproxy-session-ticket-keys-rotate
